# salt/base/files/upstart/metricbeat.conf
description "Filebeat container"
start on started logstash
stop on stopping docker

pre-start script
  /bin/sh -c 'if [ "$(/usr/bin/docker ps -aq -f name=metricbeat)" ]; then /usr/bin/docker rm metricbeat ; fi'
  exec /usr/bin/docker login -u {{ docker_user }} -p {{ docker_password }}
  exec /usr/bin/docker pull {{ docker_user }}/{{ my_metricbeat }}:latest
end script

pre-stop script
    /usr/bin/docker kill metricbeat
end script

respawn

exec /bin/sh -c '\
    docker run  -e LOGSTASH_HOST=logstash:5044 \
                --hostname=`cat /etc/hostname` \
                --name metricbeat \
                --user=root \
                --network=mynet \
                -v /var/run/docker.sock:/var/run/docker.sock \
                --mount type=bind,source=/usr/share/metricbeat/metricbeat.yml,target=/usr/share/metricbeat/metricbeat.yml \
                --mount type=bind,source=/proc,target=/hostfs/proc,readonly \
                --mount type=bind,source=/sys/fs/cgroup,target=/hostfs/sys/fs/cgroup,readonly \
                --mount type=bind,source=/,target=/hostfs,readonly \
                {{ docker_user }}/{{ my_metricbeat }}:latest -system.hostfs=/hostfs'

